# crontab
* * * * * root /tasks/photo-mover.sh >> /tasks/cron.log 2>&1
